# rsm-dev-468/rpmbuild/SOURCES/cron.d/rsm-slv-dns-udp-rtt-pl
* * * * * root sleep 11; /opt/zabbix/scripts/slv/rsm.slv.dns.udp.rtt.pl >> /var/log/zabbix/rsm.slv.err 2>&1
