# rsm-dev/rpmbuild/SOURCES/cron.d/sla-monthly-status-pl
0 15 1 * * root /opt/zabbix/scripts/slv/sla-monthly-status.pl >> /var/log/zabbix/sla-monthly-status.err 2>&1
