# rsm-qa-fix-283/rpmbuild/SOURCES/cron.d/rsm-slv-rdds-downtime-pl
* * * * * root sleep 7; /opt/zabbix/scripts/slv/rsm.slv.rdds.downtime.pl >> /var/log/zabbix/rsm.slv.err 2>&1
