# rsm-qa/rpmbuild/SOURCES/cron.d/rsm-slv-dnssec-rollweek-pl
* * * * * root sleep 3; /opt/zabbix/scripts/slv/rsm.slv.dnssec.rollweek.pl >> /var/log/zabbix/rsm.slv.err 2>&1
