# rsm-prod-qa-fix/rpmbuild/SOURCES/cron.d/rsm-slv-rdds-rollweek-pl
* * * * * root sleep 5; /opt/zabbix/scripts/slv/rsm.slv.rdds.rollweek.pl >> /var/log/zabbix/rsm.slv.err 2>&1
