# rsm-qa-fix/rpmbuild/SOURCES/cron.d/rsm-slv-dns-downtime-pl
* * * * * root sleep 6; /opt/zabbix/scripts/slv/rsm.slv.dns.downtime.pl >> /var/log/zabbix/rsm.slv.err 2>&1
