# rsm-qa-fix/rpmbuild/SOURCES/cron.d/rsm-slv-dns-rollweek-pl
* * * * * root sleep 4; /opt/zabbix/scripts/slv/rsm.slv.dns.rollweek.pl >> /var/log/zabbix/rsm.slv.err 2>&1
