# rsm-prod-qa/rpmbuild/SOURCES/cron.d/rsm-slv-dns-avail-pl
* * * * * root sleep 2; /opt/zabbix/scripts/slv/rsm.slv.dns.avail.pl >> /var/log/zabbix/rsm.slv.err 2>&1
