# rsm-qa-fix-223/rpmbuild/SOURCES/cron.d/rsm-slv-rdds-avail-pl
* * * * * root sleep 0; /opt/zabbix/scripts/slv/rsm.slv.rdds.avail.pl >> /var/log/zabbix/rsm.slv.err 2>&1
