# rsm-prod/rpmbuild/SOURCES/cron.d/rsm-probe-online-pl
* * * * * root /opt/zabbix/scripts/slv/rsm.probe.online.pl >> /var/log/zabbix/rsm.slv.err 2>&1
